# backend_api.service
[Unit]
Description=API REST backend IoT (FastAPI)
After=network.target

[Service]
Type=simple
User=daniel
WorkingDirectory=/home/daniel/repos/iot_streamlit
Environment=VIRTUAL_ENV=/home/daniel/repos/iot_streamlit/.iot_streamlit
Environment=PATH=/home/daniel/repos/iot_streamlit/.iot_streamlit/bin:/usr/bin:/bin
ExecStart=/home/daniel/repos/iot_streamlit/.iot_streamlit/bin/python -m uvicorn backend.api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
